# etc/systemd/system/update-service-ips.service
[Unit]
Description=Update service IP ranges
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-service-ips.sh

[Install]
WantedBy=multi-user.target
